# com.neoagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.neoagent</string>

    <key>ProgramArguments</key>
    <array>
        <string>__NODE_BIN__</string>
        <string>__APP_DIR__/server/index.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__APP_DIR__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>NEOAGENT_HOME</key>
        <string>__RUNTIME_HOME__</string>
    </dict>

    <!-- Auto-restart if the process exits for any reason -->
    <key>KeepAlive</key>
    <true/>

    <!-- Start immediately when the plist is loaded / on login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Throttle rapid restart loops: wait 10s before restarting -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/neoagent.log</string>

    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/neoagent.error.log</string>
</dict>
</plist>
